# provision/playbooks/templates/mailcatcher.service
[Unit]
 Description=Mailcatcher Service

[Service]
 Type=simple
 ExecStart=/usr/local/bin/mailcatcher --foreground --http-ip=0.0.0.0 --no-quit
